# overlays/RoninOS/usr/lib/systemd/system/ronin-setup.service
[Unit]
Description=Setup RoninDojo
After=oem-boot.service

[Service]
Type=simple
ExecStart=/opt/setup/ronin-setup.sh
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
