# overlays/RoninOS/usr/lib/systemd/system/ronin-setup.service
[Unit]
Description=Setup RoninDojo
After=avahi-daemon.service network-online.target
Wants=oem-boot.service network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin
ExecStart=/usr/local/sbin/ronin-setup.sh
ExecStop=/usr/local/sbin/ronin-teardown.sh
StandardOutput=file:/home/ronindojo/.logs/setup.logs

[Install]
WantedBy=multi-user.target